# test-configs/rootfs/overwrite/etc/rc.local
#!/bin/bash

echo "Executing rc.local"

echo "none" > /sys/class/leds/work-led/trigger

exit 0